# .config/supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
user = root

[program:nginx]
command = nginx -g 'daemon off;'

[program:gunicorn]
command = gunicorn -c /srv/sofastcar/.config/gunicorn.py config.wsgi.staging